# scrapeNews/web_app/install/systemd/web_app.service
[Unit]
Description=Gunicorn instance to serve web_app
After=network.target

[Service]
User=lucifer
Group=www-data
WorkingDirectory=/home/lucifer/extdrive/Documents/thisisayush/scrape/scrapeNews/web_app
Environment="PATH=/home/lucifer/extdrive/Documents/thisisayush/deploy/venv/bin"
ExecStart=/home/lucifer/extdrive/Documents/thisisayush/scrape/scrapeNews/web_app/start

[Install]
WantedBy=multi-user.target
